# web/cron-conf/crontab
# REDCap Cron Job (runs every minute)
* * * * * /usr/bin/php /app/redcap/cron.php > /dev/null
